# packer/files/fetch-fresh-ip.service
[Unit]
Description="Fetch fresh IP (new clone)"

[Service]
Type=oneshot
ExecStart=/bin/bash -c "/usr/bin/networkctl renew $(ip link show | awk -F ': ' '{print $2}' | awk '{print $1}' | grep ens)"
[Install]
WantedBy=network-online.target